# config/leximancer.conf
[program:leximancer]
command='/home/ubuntu/leximancer/current/bin/LexServer.sh console'
directory=/home/ubuntu/leximancer/current
user=ubuntu
autostart=true
autorestart=true
redirect_stderr=True
